# bar.service
[Unit]
Description=the bar service
After=network.target

[Service]
ExecStart=/usr/bin/stress --cpu 4

[Install]
WantedBy=multi-user.target

